# scripts/com.ppe-platform.news-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ppe-platform.news-collector</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/maxiaoha/Desktop/mdlooker/mdlooker/ppe-platform/scripts/news-collector.js</string>
        <string>--daily</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/ppe-news-collector.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ppe-news-collector-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/maxiaoha/Desktop/mdlooker/mdlooker/ppe-platform</string>
</dict>
</plist>
